# scripts/assistente-vi-app.service
[Unit]
Description=Assistente Vi - Amo Vidas (Next.js)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/www/vi-amovidas/assistente-vi
Environment=NODE_ENV=production
Environment=PORT=3001
ExecStart=/usr/bin/npm run start:prod
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=assistente-vi-app

[Install]
WantedBy=multi-user.target
